# collections/ansible_collections/arcde/kotami/roles/niri_lisgd/files/etc/systemd/user/niri-lisgd.service
[Unit]
Description=Niri Touch Gestures Service
After=graphical.target
Wants=graphical.target

[Service]
Type=exec
ExecStart=/usr/local/bin/niri-lisgd
Restart=on-failure
RestartSec=3
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0

# Run as the logged-in user
User=%i
Group=%i

# Required for accessing input devices
SupplementaryGroups=input

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only

[Install]
WantedBy=graphical.target
