# collections/ansible_collections/arcfg/env/roles/editor/templates/etc/profile.d/00.editor.sh.j2
#!/bin/sh

export EDITOR={{ editor_env_arcfg_executable }}
