# collections/ansible_collections/arcde/kotami/roles/niri_lisgd/files/home/user/.config/systemd/user/niri-lisgd-setup.service
[Unit]
Description=Niri LISGD Device Setup
Before=niri-lisgd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=%h/.local/bin/niri-lisgd-setup

[Install]
WantedBy=default.target
