# collections/ansible_collections/arcde/kotami/roles/niri_lisgd/files/home/user/.config/systemd/user/niri-lisgd.service
[Unit]
Description=Touch Screen Gesture Daemon
After=graphical-session.target
Wants=graphical-session.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
ExecStart=%h/.local/bin/niri-lisgd
Restart=on-failure
RestartSec=5
StartLimitAction=none

[Install]
WantedBy=graphical-session.target
